# packages/python/system/plent/system.d/plent.service
[Unit]
Description=Plent
Requires=mpd.service
After=mpd.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/plent.py

[Install]
WantedBy=multi-user.target
