# packages/graphics/rpi-fbcp/system.d/fbcp.service
[Unit]
Description=fbcp

[Service]
Restart=always
ExecStart=/usr/bin/fbcp

[Install]
WantedBy=multi-user.target
